# systemd/telegram-admin.service
[Unit]
Description=Telegram Wellness Bot - Admin API
After=network.target

[Service]
Type=simple
User=telegram-bot
WorkingDirectory=/opt/telegram_wellness_bot
Environment="PATH=/opt/telegram_wellness_bot/venv/bin"
# Bind to loopback only. The admin panel uses plaintext HTTP Basic auth and
# exposes high-risk tools -- expose it via a TLS reverse proxy (nginx), not by
# binding 0.0.0.0 directly.
ExecStart=/opt/telegram_wellness_bot/venv/bin/python -m app.interfaces.admin.server --host 127.0.0.1 --port 8110
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
