# systemd/telegram-admin.service
[Unit]
Description=Telegram Wellness Bot - Admin API
After=network.target

[Service]
Type=simple
User=telegram-bot
WorkingDirectory=/opt/telegram_wellness_bot
Environment="PATH=/opt/telegram_wellness_bot/venv/bin"
ExecStart=/opt/telegram_wellness_bot/venv/bin/python -m app.interfaces.admin.server --host 0.0.0.0 --port 8100
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
